# packaging/systemd/lemonade-nexus.service
[Unit]
Description=Lemonade-Nexus Mesh VPN Server
Documentation=https://github.com/OWNER/lemonade-nexus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lemonade-nexus
Group=lemonade-nexus
ExecStart=/usr/bin/lemonade-nexus
WorkingDirectory=/var/lib/lemonade-nexus
EnvironmentFile=-/etc/lemonade-nexus/lemonade-nexus.env
Restart=always
RestartSec=5
LimitNOFILE=65536

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/lib/lemonade-nexus
ReadOnlyPaths=/etc/lemonade-nexus

# Allow binding to privileged ports (HTTP 80/443, DNS 53) without running as
# root. CAP_NET_ADMIN is intentionally NOT granted: the boringtun dataplane runs
# fully in userspace (no kernel interface / TUN device), so the daemon needs no
# network-administration privileges.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
